# board/monome/common/rootfs-overlay/etc/systemd/system/norns-matron.service
[Unit]

[Service]
Type=simple
User=root
ExecStart=/bin/ws-wrapper ws://*:5555 /bin/matron
