# board/norns/rootfs-overlay/etc/systemd/system/norns-matron.service
[Unit]
Description=Matron

[Service]
Type=simple
# User=we
# Group=we
User=root
ExecStart=/bin/ws-wrapper ws://*:5555 /bin/matron
KillMode=process
